# system_files/usr/lib/systemd/system/goose-boot-ok.timer
[Unit]
Description=Mark goose boot as successful after boot settles
Requires=goose-boot-ok.service

[Timer]
OnBootSec=2m

[Install]
WantedBy=timers.target
