# system_files/usr/lib/systemd/system/goose-boot-ok.service
[Unit]
Description=Mark goose boot as successful
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /var/lib/goose-boot-ok
